# crontab
0 0 * * * /usr/bin/systemctl restart secapture
0 0 * * * /usr/bin/systemctl restart seconvert

