# se.service
[Unit]
Description=SolarEdge app
After=syslog.target
[Service]
ExecStart=/root/se
Restart=on-abort
[Install]
WantedBy=multi-user.target

