# scripts/crontab
0 0 * * * /bin/systemctl restart solaredge 
